# s812/etc/rc.local
#!/bin/sh

dhclient eth0 &

/etc/hdmi.sh &
