# s812/etc/rc.local
#!/bin/sh

/etc/hdmi.sh &

dhclient eth0 &

exit 0